# deploy/systemd/elliott-wave-edge-agent.service
[Unit]
Description=Elliott Wave Edge Agent (Gemini analytics)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/wave-structure-engine
ExecStart=/root/wave-structure-engine/.venv/bin/python /root/wave-structure-engine/services/edge_agent.py
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/root/wave-structure-engine/.env

[Install]
WantedBy=multi-user.target
